# APP-META/docker-config/environment/alibaba/kube-apiserver/bin/rc.local
#!/bin/bash

# called by sigma-keepalive plugin
source /etc/bashrc

#ip addr show docker0 && {
#    service staragentctl stop
#    service DragoonAgent stop
#    chkconfig --level 3 DragoonAgent off
#    chkconfig --level 3 staragentctl off
#}

chown -R admin:admin /home/admin

workDir=/home/admin/kube-apiserver
sudo -E -u admin /bin/bash $workDir/bin/start-k8s-apiserver.sh &> $workDir/log/start.log
